# packaging/deb/postinst
#!/bin/sh -e

if [ "$1" = "configure" ]; then
    if [ -z "$2" ]; then
        # initial installation
        if [ ! -f /var/lib/bee/keys/libp2p.key ]; then
            /usr/bin/bee init --config /etc/bee/bee.yaml >/dev/null
            chown -R bee:bee /var/lib/bee
        fi
        echo "
Logs:   journalctl -f -u bee.service
Config: /etc/bee/bee.yaml

Bee requires a Gnosis Chain RPC endpoint to function. By default this is expected to be found at ws://localhost:8546.

Please see https://docs.ethswarm.org/docs/installation/install for more details on how to configure your node.

After you finish configuration run 'sudo bee-get-addr' and fund your node with XDAI, and also XBZZ if so desired.
        "
    fi

    deb-systemd-helper unmask bee.service >/dev/null || true

    if deb-systemd-helper --quiet was-enabled bee.service; then
        deb-systemd-helper enable bee.service >/dev/null || true
    else
        deb-systemd-helper update-state bee.service >/dev/null || true
    fi
    if [ -d /run/systemd/system ]; then
        systemctl --system daemon-reload >/dev/null || true
        # upgrade
        if [ ! -z "$2" ]; then
            deb-systemd-invoke start bee.service >/dev/null || true
        fi
    fi
fi
