# packaging/deb/postinst
#!/bin/sh -e

deb-systemd-helper unmask bee.service >/dev/null || true

if deb-systemd-helper --quiet was-enabled bee.service; then
    deb-systemd-helper enable bee.service >/dev/null || true
else
    deb-systemd-helper update-state bee.service >/dev/null || true
fi
if [ -d /run/systemd/system ]; then
    systemctl --system daemon-reload >/dev/null || true
    deb-systemd-invoke start bee.service >/dev/null || true
fi