# packaging/bee.service
[Unit]
Description=Ethereum Swarm node
Documentation=https://swarm-gateways.net/bzz:/docs.swarm.eth/
After=network.target

[Service]
NoNewPrivileges=true
User=bee
Group=bee
StateDirectory=bee
StateDirectoryMode=0750
ConfigurationDirectory=bee
ExecStart=/usr/bin/bee start --config=/etc/bee/bee.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
