# packaging/bee.service
[Unit]
Description=Ethereum Swarm node
Documentation=https://swarm-gateways.net/bzz:/docs.swarm.eth/
After=network.target

[Service]
<<<<<<< HEAD
EnvironmentFile=/etc/default/bee
=======
>>>>>>> remove cap_bind and ulimit
NoNewPrivileges=true
User=bee
Group=bee
ExecStart=/usr/bin/bee start --config=/etc/bee/bee.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
